# deploy/systemd/mailflow-monitor.timer
[Unit]
Description=Run Mailflow Monitor every five minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Unit=mailflow-monitor.service

[Install]
WantedBy=timers.target

