# deploy/systemd/mailflow-monitor.service
[Unit]
Description=Mailflow Monitor check

[Service]
Type=oneshot
WorkingDirectory=%h/.config/mailflow-monitor
ExecStart=%h/.local/share/mailflow-monitor/venv/bin/mailflow-monitor check --config %h/.config/mailflow-monitor/config.toml
